# queue.work.ini
[program:queue-work]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/youtube-audio/artisan queue:work --queue=default --tries=1 --timeout=180
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel.queue.work.log
